# apps/background-removal/background-removal-jpyc.service
[Unit]
Description=JPYC x402 Payment Gateway for Background Removal
After=network.target background-removal.service
Requires=background-removal.service

[Service]
Type=simple
WorkingDirectory=/home/kojima/work/url2ai/apps/background-removal
EnvironmentFile=-/home/kojima/work/url2ai/apps/background-removal/.env.jpyc
Environment=PATH=/home/kojima/.local/node-v22.15.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/kojima/.local/node-v22.15.0/bin/node server-jpyc.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
