# src/polymarket_worker.service
[Unit]
Description=URL2AI Polymarket RQ Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/kojima/work/url2ai
ExecStart=/home/kojima/work/url2ai/scripts/run_polymarket_rq_worker.sh
Restart=always
RestartSec=20
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:/home/kojima/rqdb4ai_polymarket_worker.log
StandardError=append:/home/kojima/rqdb4ai_polymarket_worker.log

[Install]
WantedBy=multi-user.target
