# apps/llm-gateway/fxbrain-jpyc.service
[Unit]
Description=Kurage FX Brain JPYC x402 Payment Gateway (Gemma 4 12B)
After=network.target

[Service]
WorkingDirectory=/home/kojima/work/url2ai/apps/llm-gateway
ExecStart=/home/kojima/.local/node-v22.15.0/bin/node server-jpyc-fxbrain.js
Restart=always
RestartSec=5
EnvironmentFile=-/home/kojima/work/url2ai/apps/llm-gateway/.env.jpyc
EnvironmentFile=-/home/kojima/work/kfxbrain/.env
Environment=PATH=/home/kojima/.local/node-v22.15.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PORT_JPYC_FXBRAIN=18327
Environment=FXBRAIN_URL=http://127.0.0.1:18326

[Install]
WantedBy=multi-user.target
